# deploy/km-api.service.example
# systemd unit — API only (nginx serves frontend)
#
# Install:
#   sudo cp deploy/km-api.service.example /etc/systemd/system/km-api.service
#   sudo nano /etc/systemd/system/km-api.service   # fix User= and WorkingDirectory=
#   sudo systemctl daemon-reload
#   sudo systemctl enable km-api
#   sudo systemctl start km-api
#   sudo systemctl status km-api
#   journalctl -u km-api -f

[Unit]
Description=KM Dynasty API (Express + Postgres)
After=network.target postgresql.service

[Service]
Type=simple
User=YOUR_LINUX_USER
WorkingDirectory=/home/YOUR_LINUX_USER/TIKTOK-WEBSITE/server
Environment=NODE_ENV=production
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
